# com.valhalla-farms.omnifocus-motion-healthcheck.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.valhalla-farms.omnifocus-motion-healthcheck</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/valhalla-farms/scripts/venv/bin/python3</string>
        <string>/Users/valhalla-farms/scripts/health_check.py</string>
    </array>

    <!-- Run every 30 minutes -->
    <key>StartInterval</key>
    <integer>1800</integer>

    <!-- Run at load -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Standard output log -->
    <key>StandardOutPath</key>
    <string>/Users/valhalla-farms/Library/Logs/OmniFocusMotionSync/healthcheck.log</string>

    <!-- Error log -->
    <key>StandardErrorPath</key>
    <string>/Users/valhalla-farms/Library/Logs/OmniFocusMotionSync/healthcheck_error.log</string>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/valhalla-farms/scripts</string>
</dict>
</plist>
