# com.valhalla-farms.omnifocus-motion-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.valhalla-farms.omnifocus-motion-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/valhalla-farms/scripts/run_sync.sh</string>
    </array>

    <!-- Run every 15 minutes -->
    <key>StartInterval</key>
    <integer>900</integer>

    <!-- Run at load (system startup) -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Standard output log -->
    <key>StandardOutPath</key>
    <string>/Users/valhalla-farms/Library/Logs/OmniFocusMotionSync/sync.log</string>

    <!-- Error log -->
    <key>StandardErrorPath</key>
    <string>/Users/valhalla-farms/Library/Logs/OmniFocusMotionSync/error.log</string>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/valhalla-farms/scripts</string>
</dict>
</plist>
