# com.yourusername.omnifocus-motion-sync.example.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.yourusername.omnifocus-motion-sync</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Update these paths to match your installation -->
        <string>/path/to/motion-sync/venv/bin/python3</string>
        <string>/path/to/motion-sync/sync_of_to_motion.py</string>
        <string>--config</string>
        <string>/path/to/motion-sync/config.json</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>MOTION_API_KEY</key>
        <string>YOUR_API_KEY_HERE</string>
    </dict>

    <!-- Run every 15 minutes -->
    <key>StartInterval</key>
    <integer>900</integer>

    <!-- Run at load (system startup) -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Standard output log (LaunchAgents don't expand ~) -->
    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/Library/Logs/OmniFocusMotionSync/sync.log</string>

    <!-- Error log -->
    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/Library/Logs/OmniFocusMotionSync/error.log</string>

    <!-- Working directory - update to your installation path -->
    <key>WorkingDirectory</key>
    <string>/path/to/motion-sync</string>

    <!-- See also: com.yourusername.omnifocus-motion-healthcheck.example.plist for health monitoring -->
</dict>
</plist>
